# extras/openbsd/fwknop-2.0.4/pkg/fwknopd.rc
#!/bin/sh
#
# $OpenBSD$

daemon="${TRUEPREFIX}/sbin/fwknopd"

. /etc/rc.d/rc.subr

rc_cmd $1
